# toolbox/certbot/certbot_renewals.cron
# Cron jobs for renewing SSL/TLS certificates using Certbot
0 7,19 * * 1-5 root /srv/dts-services-haproxy/toolbox/certbot/renew_ssl_tls_certbot_certificate.sh -d atd-postgrest.austinmobility.io -t /srv/dts-services-haproxy/ssl -c haproxy
0 7,19 * * 1-5 root /srv/dts-services-haproxy/toolbox/certbot/renew_ssl_tls_certbot_certificate.sh -d equitytool.austinmobility.io -t /srv/dts-services-haproxy/ssl -c haproxy
